# roles/ansible-role-volumio-bluetooth/files/lib/systemd/system/bluealsa.service
[Unit]
Description=BluezAlsa proxy
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
User=root
Group=audio
ExecStart=/usr/bin/bluealsa --profile=a2dp-sink

[Install]
WantedBy=multi-user.target
