# roles/ansible-role-volumio-bluetooth/files/etc/udev/rules.d/99-input.rules
KERNEL=="input[0-9]*", RUN+="/home/volumio/a2dp-autoconnect"
